# dockerSrc/conf/crontab
# RNMS
*/5 * * * * /RNMS/bin/checkStatus.sh >> /RNMS/log/checkStatus.log 2>&1
*/5 * * * * /RNMS/bin/checkPerf.sh 5\ >> /RNMS/log/checkPerf.log 2>&1
